# etc/update-motd.d/99-intelcom
#!/bin/bash

echo -e '\033[34m#######################################################
#\033[33m   PBX Intelcom - VoIP АТС Производства "Интелком"   \033[34m#
#\033[33m            tel.atsip.ru  Тел.(3532)XX-XX-XX         \033[34m#
#######################################################\033[0m'

{
echo "              _       _                   "
echo "         __ _| |_ ___(_)_ __   _ __ _   _ "
echo "        / _\` | __/ __| | '_ \ | '__| | | |"
echo "       | (_| | |_\__ \ | |_) || |  | |_| |"
echo "        \__,_|\__|___/_| .__(_)_|   \__,_|"
echo "                       |_|                "
}

